# cookbooks/reveal.js/files/etc/systemd/system/reveal.js.service
[Service]
WorkingDirectory=/opt/reveal.js
ExecStart=/usr/bin/grunt serve
Restart=always

[Install]
WantedBy=multi-user.target
