# cookbooks/letsencrypt/templates/etc/cron.d/letsencrypt.erb
MAILTO="<%= @mailto %>"
5 3 * * * root /opt/letsencrypt/letsencrypt.sh -c && <%= @services.map {|service| "systemctl reload #{service}" }.join(' && ') %>
