# cookbooks/plant/files/etc/systemd/system/plant.service
[Unit]
After=plant.target

[Service]
ExecStart=/opt/plant/bin/plant start
Restart=always

[Install]
WantedBy=multi-user.target
