# cookbooks/kibana/files/etc/systemd/system/kibana.service
[Unit]
After=elasticserch.target

[Service]
ExecStart=/opt/kibana/bin/kibana
Restart=always
User=kibana

[Install]
WantedBy=multi-user.target
